# server/systemd/simplepresent.service

[Unit]
Description=SimplePresent sync server
After=network.target

[Service]
Type=simple
User=simplepresent
Group=simplepresent
WorkingDirectory=/var/lib/simplepresent
ExecStart=/usr/local/bin/simplepresent -config /etc/simplepresent/config.json
Restart=on-failure
RestartSec=5
Environment=CONFIG=/etc/simplepresent/config.json
RuntimeDirectory=simplepresent
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
